# systemd/claude-resume-later.timer
[Unit]
Description=Run claude-resume-later every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
RandomizedDelaySec=15s

[Install]
WantedBy=timers.target
